# packages/bsp/jethub/jethubj80/05-jethub_set_eth_leds.rules
SUBSYSTEM=="net",KERNEL=="eth0",RUN+="/usr/sbin/jethub_set-eth_leds"
