# packages/bsp/meson-s4t7/fan.service
[Unit]
Description=Automatic Fan Control Service

[Service]
Type=oneshot
RemainAfterExit=Yes
ExecStart=/bin/sh -c "echo 1 > /sys/class/fan/mode"
ExecStart=/bin/sh -c "echo 1 > /sys/class/fan/enable"

[Install]
WantedBy=multi-user.target
