# packages/bsp/ti/weston/weston.socket
[Unit]
Description=Weston socket
RequiresMountsFor=/run

[Socket]
ListenStream=/run/wayland-0
SocketMode=0775
# SocketUser=weston
# SocketGroup=wayland
RemoveOnStop=yes

[Install]
WantedBy=sockets.target

