# config/packages/99-board/dedicated/olimex-som204-a20/overlay/lib/systemd/system/olinuxino-lcd.service
[Unit]
Description=LCD-OLinuXino-15.6FHD source clock tune
After=basic.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/devmem 0x01c20118 w 0xc3000000

[Install]
WantedBy=multi-user.target
