# config/packages/03-common/dedicated/02-armbian-desktop-xfce/overlay/etc/profile.d/sound.sh
if [ -z "$GTK_MODULES" ] ; then
	GTK_MODULES="canberra-gtk-module"
else
	GTK_MODULES="$GTK_MODULES:canberra-gtk-module"
fi

export GTK_MODULES
