# packages/bsp/xpressreal-t3/load-rtk-modules.service
[Unit]
Description=Load RTK Kernel Modules in a Specific Order
DefaultDependencies=no
Conflicts=shutdown.target
Before=sysinit.target shutdown.target
ConditionCapability=CAP_SYS_MODULE

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/load-rtk-modules.sh
TimeoutSec=90s

[Install]
WantedBy=multi-user.target

