# extensions/network/config-networkd/systemd/system/apt-daily-upgrade.service.d/armbian-exec-correct-service.conf
[Service]
ExecStartPre=
ExecStartPre=/lib/systemd/systemd-networkd-wait-online --any --timeout=60
