# packages/bsp/rk3399/sprd-bluetooth.service
[Unit]
Description=Spreadtrum (sprd) bluetooth support
Before=bluetooth.service

[Service]
Type=simple
RemainAfterExit=Yes
ExecStartPre=/usr/sbin/rfkill unblock all
ExecStart=/usr/bin/hciattach_opi -n -s 1500000 /dev/ttyBT0 sprd

[Install]
WantedBy=multi-user.target
