# packages/bsp/common/etc/udev/rules.d/10-wifi-disable-powermanagement.rules
KERNEL=="wlan*", ACTION=="add", RUN+="/sbin/iw dev %k set power_save off"
