# packages/bsp/common/etc/udev/rules.d/10-wifi-disable-powermanagement.rules
SUBSYSTEM=="net", ACTION!="remove", ENV{DEVTYPE}=="wlan", RUN+="/usr/sbin/iw dev $name set power_save off"
