# packages/bsp/common/etc/cron.d/armbian-truncate-logs
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

*/15 * * * * root /usr/lib/armbian/armbian-truncate-logs
@reboot root /usr/lib/armbian/armbian-truncate-logs
