# packages/bsp/common/etc/cron.d/armbian-updates
@reboot root /usr/lib/armbian/apt-updates
@daily root /usr/lib/armbian/apt-updates
