# packages/bsp/common/etc/cron.d/armbian-truncate-logs
PATH=/usr/bin:/bin:/usr/sbin:/sbin

*/15 * * * * root /usr/lib/armbian/armbian-truncate-logs

