# packages/bsp/common/etc/cron.d/armbian-truncate-logs
*/15 * * * root /usr/lib/armbian/armbian-truncate-logs