# packages/bsp/common/etc/NetworkManager/dispatcher.d/30-offload
#!/usr/bin/env bash
# With offloading enabled, Samba becomes completely useless - transfers to a Windows PC freezes randomly
# Disabling offloading should solve this - all RK3399

IFACE=$1
EVENT=$2

if [[ "$EVENT" != "up" ]]; then
        exit 1
fi

ETHTOOL=/sbin/ethtool

[ -f $ETHTOOL ] || exit 2
[ -d /sys/devices/platform/fe300000.ethernet/net/$IFACE ] || exit 3

$ETHTOOL -K $IFACE rx off tx off
exit 0
