# packages/bsp/radxa3308/etc/udev/rules.d/05-rockpis-MACaddress.rules
#If a network interface is being assigned a different MAC on each boot, or the
#MAC address is based on the disk image, rather than being unique to hardware,
#enable the corresponding line below to derive that interface's MAC from
#the RK3308 SOC's unique serial number.
#All the two digit hex code prefixes passed to fixEtherAddr should be unique
#and chosen such that (n-2)%4 == 0

#Otherwise, the RockPI-S built-in ethernet MAC is derived from rootfs image
KERNEL=="end0", ACTION=="add" RUN+="fixEtherAddr %k 06"

#Otherwise, the RockPI-S wifi chip MACs will be different on each boot
KERNEL=="wlan0", ACTION=="add" RUN+="fixEtherAddr %k 0a"
KERNEL=="p2p0", ACTION=="add" RUN+="fixEtherAddr %k 0e"
