# userpatches/overlay/rc.local
#!/bin/bash

echo "[rc.local] - start at $(date '+%Y-%m-%d %H:%M:%S')"

# Example: Save the startup date to a file in /tmp
date > /tmp/rc_local_started.txt

echo "[rc.local] - exit 0 at $(date '+%Y-%m-%d %H:%M:%S')"
exit 0