# packages/bsp/rk322x/rtk-bluetooth.service
[Unit]
Description=Realtek H5 protocol bluetooth support 
Before=bluetooth.service

[Service]
Type=forking
RemainAfterExit=Yes
ExecStart=/usr/bin/rtk_hciattach -n -s 115200 /dev/ttyS1 rtk_h5 1500000

[Install]
WantedBy=multi-user.target
