# packages/bsp/ayn-odin2/bt-fixed-mac.service
[Unit]
Description=Bluetooth Fixed Mac Address
After=bluetooth.service
Wants=bluetooth.service
[Service]
Type=oneshot
RemainAfterExit=yes
SyslogIdentifier=bt-fied-mac
ExecStart=/usr/local/bin/bt-fied-mac.sh
[Install]
WantedBy=bluetooth.target
