# config/cli/bsp/common/lib/systemd/system/armbian-disable-autologin.service
[Unit]
Description=Disable automated desktop login

[Service]
Type=simple
ExecStart=rm -f /etc/lightdm/lightdm.conf.d/22-armbian-autologin.conf
TimeoutStopSec= 10
ExecStop=systemctl disable armbian-disable-autologin.timer
