# packages/bsp/easepi/bt-hciattach.service
[Unit]
Description=EasePi Bluetooth HCIAttach
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStart=/usr/local/sbin/bt-hciattach.sh

[Install]
WantedBy=multi-user.target
