# packages/bsp/easepi/ir.service
[Unit]
Description=EasePi IR Keymap Configuration
After=systemd-modules-load.service
After=sysinit.target

[Service]
Type=oneshot
ExecStart=/usr/local/ir/ir.sh --quick
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
