# packages/bsp/arduino/armbian-resize-filesystem-qcom.service
[Unit]
Description=Expand rootfs to fill eMMC
DefaultDependencies=no
Before=armbian-firstrun.service
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/lib/armbian/armbian-resize-filesystem-qcom
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
