# packages/bsp/common/etc/cron.weekly/armbian-quotes
#!/bin/bash

curl -L -s -o "/etc/update-motd.d/quotes.txt" "https://dl.armbian.com/quotes.txt"
