# config/packages/03-common/dedicated/02-armbian-desktop-xfce/overlay/etc/skel/.config/autostart/loginsound.desktop
[Desktop Entry]
Encoding=UTF-8
Version=0.9.4
Type=Application
Name=Login sound
Comment=
Exec=canberra-gtk-play -f /usr/share/sounds/Smooth/stereo/desktop-login.oga
OnlyShowIn=XFCE;
StartupNotify=false
Terminal=false
Hidden=false
