# scripts/update-motd.d/10-header
#!/bin/bash
export TERM=linux
. /etc/os-release
KERNELID=$(uname -r)
([[ -f /etc/armbian-release ]] && source /etc/armbian-release || read ID </run/machine.id) 2>/dev/null
read HOSTNAME </etc/hostname
clear
toilet -f standard -F metal ${ID:-$HOSTNAME}
echo ""
printf 'Welcome to \e[0;91mARMBIAN\x1B[0m %s %s\n' "$PRETTY_NAME $KERNELID"
echo ""
