# scripts/update-motd.d/10-header
#!/bin/bash
export TERM=linux
export VDPAU_DRIVER=sunxi 
. /etc/os-release
KERNELID=$(uname -r) 
clear
toilet -f standard -F metal  $(cat /var/run/machine.id)
echo ""
printf 'Welcome to \e[0;91mARMBIAN\x1B[0m %s %s\n' "$PRETTY_NAME $KERNELID"
echo ""
