# scripts/update-motd.d/41-armbian-config
#!/bin/bash

echo -e "[\e[31m General system configuration\e[0m: \e[1marmbian-config\e[0m ]\n"