# scripts/update-motd.d/99-point-to-faq
#!/bin/bash

. /etc/armbian-release

if [ "$LINUXFAMILY" = "sun8i" ]; then
	printf "\nNew to Armbian? Check the Armbian H3 Mini FAQ first:\n\e[0;92m https://github.com/igorpecovnik/lib/blob/master/documentation/H3_mini_faq.md \x1B[0m\n"
fi

(sleep 1 && rm "$0") &
