# scripts/update-motd.d/99-point-to-faq
#!/bin/bash

. /etc/armbian-release

printf "New to Armbian? Check the documentation first: \e[0;92mdocs.armbian.com\x1B[0m\n\n"

(sleep 1 && rm "$0") &
