# scripts/update-motd.d/99-point-to-faq
#!/bin/bash

. /etc/armbian-release

printf "New to Armbian? Check the documentation first: \e[0;92mhttps://docs.armbian.com\x1B[0m\n\n"

if [ -f "/usr/bin/armbian-config" ]; then
	printf "General system configuration: \e[1marmbian-config\e[0m\n\n"
fi

(sleep 1 && rm "$0") &