# packages/bsp/oneplus-kebab/qca6390-wlan-enable.service
[Unit]
Description=QCA6390 WLAN enable (GPIO)
After=local-fs.target
Before=network-pre.target

[Service]
Type=simple
ExecStart=/usr/local/bin/qca6390-wlan-enable.sh

[Install]
WantedBy=multi-user.target
