# config/packages/03-common/99-armbian-bsp/overlay/lib/systemd/system/armbian-hardware-monitor.service
# Armbian hardware monitor service
# Detects and sets various features on boot
# This service may block the boot process for up to 2 minutes

[Unit]
Description=Armbian hardware monitoring
Before=basic.target
After=sysinit.target local-fs.target
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/lib/armbian/armbian-hardware-monitor start
ExecStop=/usr/lib/armbian/armbian-hardware-monitor stop
RemainAfterExit=yes
TimeoutStartSec=2min

[Install]
WantedBy=basic.target
