# packages/bsp/fixMACaddress/etc/udev/rules.d/05-fixMACaddress.rules
#If a network interface is being assigned a new, different address on each boot,
#enable the corresponding line below to derive its MAC addr from UUID of rootfs
#Beware that all the two digit hex code prefixes must be unique!

#KERNEL=="eth0", ACTION=="add" RUN+="fixEtherAddr %k 02"
#KERNEL=="eth1", ACTION=="add" RUN+="fixEtherAddr %k 06"
#KERNEL=="wlan0", ACTION=="add" RUN+="fixEtherAddr %k 0a"
