# config/optional/boards/anbernic-rg-vita-pro/_packages/bsp-cli/etc/systemd/system/vita-pad2key.service
[Unit]
Description=RG Vita Pro gamepad-to-keyboard remap
[Service]
ExecStart=/usr/bin/vita-pad2key
Restart=always
RestartSec=2
[Install]
WantedBy=multi-user.target
