# config/optional/boards/anbernic-rg-vita-pro/_packages/bsp-cli/etc/systemd/system/vita-jack-switch.service
[Unit]
Description=Headphone jack to internal-speaker auto switch (RG Vita Pro)
After=sound.target

[Service]
ExecStart=/usr/bin/vita-jack-switch
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
