# config/optional/boards/anbernic-rg-vita-pro/_packages/bsp-cli/etc/systemd/system/vita-offcharge.service
[Unit]
Description=Off-charge mode: battery screen on plug-in boot (RG Vita Pro)
DefaultDependencies=no
After=local-fs.target systemd-udevd.service
Before=display-manager.service getty@tty1.service graphical.target multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/bin/vita-offcharge
TimeoutStartSec=infinity

[Install]
WantedBy=multi-user.target
