# config/optional/boards/anbernic-rg-vita-pro/_packages/bsp-cli/etc/systemd/system/vita-offcharge-reboot-marker.service
[Unit]
Description=Mark next boot as warm reboot (skip off-charge screen)
DefaultDependencies=no
Before=systemd-reboot.service umount.target
RequiresMountsFor=/var/lib

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'mkdir -p /var/lib/vita-offcharge && touch /var/lib/vita-offcharge/warm-reboot'

[Install]
WantedBy=reboot.target
