# config/packages/99-board/dedicated/a5xmaxplus/overlay/etc/NetworkManager/dispatcher.d/10-skykirin-network
#!/bin/bash

if [[ $1 == "eth0" ]]
then
  if [[ $2 == "up" ]]
  then
    /usr/bin/skykirin-ht1628 -lon
  elif [[ $2 == "down" ]]
  then
    /usr/bin/skykirin-ht1628 -loff
  fi
elif [[ $1 == "wlan0" ]] || [[ $1 == "wlan1" ]]
then
  if [[ $2 == "up" ]]
  then
    /usr/bin/skykirin-ht1628 -wfull
  elif [[ $2 == "down" ]]
  then
    /usr/bin/skykirin-ht1628 -woff
  fi
fi