# packages/bsp/rk3328/z28pro/z28pro-bluetooth.service
[Unit]
Description=Realtek H5 bluetooth support
Before=bluetooth.service

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/local/bin/start_bt.sh
Restart=no

[Install]
WantedBy=multi-user.target
