# config/optional/boards/anbernic-rg-ds/_packages/bsp-cli/etc/xdg/autostart/rgds-screen.desktop
[Desktop Entry]
Type=Application
Name=Anbernic RG DS dual-screen layout
Comment=Stack the two DSI panels and map each touchscreen (X11 only)
Exec=/usr/bin/rgds-screen
NoDisplay=true
X-GNOME-Autostart-enabled=true
