# packages/bsp/common/etc/cron.daily/armbian-ram-logging
#!/bin/sh
# Only run on systems where logrotate is a cron job
if [ -f /lib/systemd/system/logrotate.timer ]; then
    exit 0
fi
/usr/lib/armbian/armbian-ramlog write >/dev/null 2>&1
