# packages/bsp/common/etc/cron.daily/armbian-quotes
#!/bin/bash

curl --max-time 2 -L -s -o "/tmp/quotes.txt" "https://github.armbian.com/quotes.txt"
