# packages/bsp/common/etc/cron.daily/log2ram
#! /bin/sh

/usr/sbin/log2ram write >/dev/null 2>&1
