# packages/bsp/common/etc/profile.d/armbian-lang.sh
# This should fix
# perl: warning: Setting locale failed.

export LC_ALL=$LANG
