# config/cli/bsp/mt7623/mt7623-wifi.service
[Unit]
Description=MT7623 Wireless
DefaultDependencies=no
Wants=rsyslog.service systemd-journald.service
Before=syslog.target sysinit.target
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mt7623-wifi.bash
RemainAfterExit=no

[Install]
WantedBy=sysinit.target
