# packages/bsp/bananapi/rtk_hciattach/bthelper.service
[Unit]
Description=Bluetooth Helper
After=bluetooth.service

[Service]
Type=simple
ExecStartPre=/usr/sbin/rfkill unblock all
ExecStart=/usr/sbin/rtk_hciattach -n -s 115200 ttyS1 rtk_h5 >/dev/null 2>&1
TimeoutSec=0
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
