# packages/bsp/gateway-dk/gateway-dk-leds.service
[Unit]
Description=Mono Gateway DK boot LED indicator
After=systemd-modules-load.service systemd-udev-settle.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/gateway-dk-leds.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
