# packages/bsp/gateway-dk/generate-fancontrol-conf.service
[Unit]
Description=Generate fancontrol configuration for Mono Gateway DK
Before=fancontrol.service
After=systemd-udev-settle.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/generate-fancontrol-conf
RemainAfterExit=yes

[Install]
WantedBy=fancontrol.service
