# packages/bsp/gateway-dk/generate-fancontrol-conf.service
[Unit]
Description=Generate fancontrol configuration for Mono Gateway DK
Before=fancontrol.service
DefaultDependencies=no
After=sys-subsystem-hwmon-devices-hwmon0.device

[Service]
Type=oneshot
ExecStart=/usr/local/bin/generate-fancontrol-conf
RemainAfterExit=yes

[Install]
WantedBy=fancontrol.service
