# packages/bsp/common/etc/update-motd.d/10-armbian-header
#!/bin/bash
#
# Copyright (c) Authors: http://www.armbian.com/authors
#
# This file is licensed under the terms of the GNU General Public
# License version 2. This program is licensed "as is" without any
# warranty of any kind, whether express or implied.

# DO NOT EDIT THIS FILE but add config options to /etc/default/armbian-motd
# any changes will be lost on board support package update

THIS_SCRIPT="header"
MOTD_DISABLE=""

[[ -f /etc/default/armbian-motd ]] && . /etc/default/armbian-motd

for f in $MOTD_DISABLE; do
	[[ $f == $THIS_SCRIPT ]] && exit 0
done

. /etc/os-release
. /etc/armbian-release

KERNELID=$(uname -r)

PRETTY_NAME=$(sed 's/.*(\(.*\)).*/\1/' <<< "$PRETTY_NAME" | awk '{print $1}' | awk '{print tolower($0)}')
[[ "$IMAGE_TYPE" != "nightly" && "$VERSION" == 20* ]] && VERSION=$(echo $VERSION | cut -f1,2 -d'.')

if [ $(echo $BOARD_NAME | wc -c) -ge 18 ]; then
	TERM=linux toilet -f standard -F metal $(echo $BOARD_NAME | sed 's/Orange Pi/OPi/' | sed 's/Nanopi/NPi/')
else
	TERM=linux toilet -f standard -F metal $BOARD_NAME
fi

echo -e "\nWelcome to \e[0;91mArmbian $PRETTY_NAME $VERSION\x1B[0m $IMAGE_TYPE $KERNELID"
