# packages/bsp/common/etc/update-motd.d/10-armbian-header
#!/bin/bash
#
# Copyright (c) Authors: http://www.armbian.com/authors
#
# This file is licensed under the terms of the GNU General Public
# License version 2. This program is licensed "as is" without any
# warranty of any kind, whether express or implied.

# DO NOT EDIT THIS FILE but add config options to /etc/default/armbian-motd
# any changes will be lost on board support package update

THIS_SCRIPT="header"
MOTD_DISABLE=""

[[ -f /etc/default/armbian-motd ]] && . /etc/default/armbian-motd

for f in $MOTD_DISABLE; do
		[[ $f == $THIS_SCRIPT ]] && exit 0
done

. /etc/os-release

# extract code name only
PRETTY_NAME=$(sed 's/.*(\(.*\)).*/\1/' <<< "$PRETTY_NAME" | awk '{print $1}' | sed 's/.*/\u&/')

# read kernel version
KERNELID=$(uname -r)

# read image data
. /etc/armbian-image-release
IMAGE_TYPE_INITIAL=$IMAGE_TYPE

# read updated data
. /etc/armbian-release

# show commit info if exists otherwise version
if [[ -n $BUILD_REPOSITORY_COMMIT ]]; then
		SHOW_VERSION="\e[0;91m"$BUILD_REPOSITORY_COMMIT"\x1B[0m"
else
		SHOW_VERSION="\e[0;92m$VERSION\x1B[0m"
fi

# stable versions
if [[ $IMAGE_TYPE_INITIAL == $IMAGE_TYPE && $IMAGE_TYPE == stable && $BUILD_REPOSITORY_URL == "https://github.com/armbian/build" ]]; then
		if [[ "$VERSION" == 20* ]]; then
		# show only prefix
		SHOW_VERSION="\e[0;92m"$(echo $VERSION | cut -f1,2 -d'.')"\x1B[0m"
		else
		SHOW_VERSION="\e[0;92m$VERSION\x1B[0m"
		fi
fi


if [ $(echo $BOARD_NAME | wc -c) -ge 18 ]; then
		TERM=linux toilet -f standard -F metal $(echo $BOARD_NAME | sed 's/Orange Pi/OPi/' | sed 's/Nanopi/NPi/')
else
		TERM=linux toilet -f standard -F metal $BOARD_NAME
fi

echo -e "Welcome to \e[0;91mArmbian\x1B[0m - $PRETTY_NAME\x1B[0m ($SHOW_VERSION) Linux $KERNELID\n"