# packages/bsp/common/etc/update-motd.d/41-armbian-config
#!/bin/bash

# DO NOT EDIT THIS FILE
# any changes will be lost on board support package update

THIS_SCRIPT="armbian-config"
MOTD_DISABLE=""

[[ -f /etc/default/armbian-motd ]] && . /etc/default/armbian-motd

for f in $MOTD_DISABLE; do
	[[ $f == $THIS_SCRIPT ]] && exit 0
done

[[ $(( $RANDOM % 2 )) == 0 && -f /usr/sbin/armbian-config ]] && echo -e "[\e[31m General system configuration\e[0m: \e[1marmbian-config\e[0m ]\n"

exit 0