# packages/bsp/common/etc/update-motd.d/41-armbian-config
#!/bin/bash
#
# Copyright (c) Authors: https://www.armbian.com/authors
#
# This file is licensed under the terms of the GNU General Public
# License version 2. This program is licensed "as is" without any
# warranty of any kind, whether express or implied.

# DO NOT EDIT THIS FILE but add config options to /etc/default/armbian-motd
# any changes will be lost on board support package update

THIS_SCRIPT="config"
MOTD_DISABLE=""

[[ -f /etc/default/armbian-motd ]] && . /etc/default/armbian-motd

for f in $MOTD_DISABLE; do
	[[ $f == $THIS_SCRIPT ]] && exit 0
done

if [[ $(( $RANDOM % 2 )) == 0 ]]; then
	if [[ -f /usr/sbin/armbian-config ]]; then
		echo -e "[\e[31m General system configuration (beta)\e[0m: \e[1marmbian-config\e[0m ]\n"
	else
		echo -e "[\e[31m Menu-driven system configuration (beta)\e[0m: \e[1msudo apt update && sudo apt install armbian-config\e[0m ]\n"
	fi
fi

exit 0

