# packages/bsp/common/etc/update-motd.d/80-update-htop
#!/bin/bash
ACTIVE_DEVICE=$(nmcli device status | grep ethernet | awk '{ print $1 }' | grep -v lo | head -1)
if [[ -f "${HOME}"/.config/htop/htoprc ]]; then
sed -i "s/^eth0_alias=.*/eth0_alias=$ACTIVE_DEVICE/" "${HOME}"/.config/htop/htoprc
fi

