# packages/bsp/rockchip/ap6330-bluetooth.service
[Unit]
Description=AP6330 bluetooth support
Before=bluetooth.service

[Service]
Type=forking
RemainAfterExit=Yes
ExecStartPre=/usr/bin/brcm_patchram_plus --patchram /lib/firmware/ap6330/brcm/bcm40183b2.hcd --enable_hci --bd_addr 14:72:23:94:a5:f6 --no2bytes --tosleep 1000 /dev/ttyS0
ExecStart=/usr/bin/hciattach /dev/ttyS0 bcm43xx 1500000

[Install]
WantedBy=multi-user.target
