# packages/bsp/rockchip/ap6330-bluetooth.service
[Unit]
Description=AP6330 bluetooth support
Before=bluetooth.service

[Service]
Type=forking
RemainAfterExit=Yes
ExecStart=/usr/bin/hciattach /dev/ttyS0 bcm43xx 1500000

[Install]
WantedBy=multi-user.target
