# packages/bsp/rockchip/tinker-bluetooth-restart.service
[Unit]
Description=Bluetooth restart walkaround
After=rc-local.service

[Service]
ExecStartPre=/bin/sleep 20
ExecStart=/bin/systemctl restart tinker-bluetooth.service
Type=idle

[Install]
WantedBy=rc-local.service