# packages/bsp/xiaomi-umi/usbgadget-rndis.service
[Unit]
Description=USB Gadget RNDIS Network
After=modprobe@usbfunc:rndis.service NetworkManager.service systemd-networkd.service
Before=dnsmasq.service
Wants=modprobe@usbfunc:rndis.service
[Service]
Type=oneshot
RemainAfterExit=yes
SyslogIdentifier=gadget-func-rndis
WorkingDirectory=/sys/kernel/config/usb_gadget
ExecStart=/usr/local/bin/setup-usbgadget-network.sh
ExecStop=/usr/local/bin/remove-usbgadget-network.sh
[Install]
WantedBy=usb-gadget.target
