# packages/bsp/aic8800/aic-bluetooth.service
[Unit]
Description=AIC8800 Bluetooth
After=bluetooth.service bluetooth.target
Requires=bluetooth.service bluetooth.target

[Service]
Type=idle
ExecStartPre=/usr/sbin/rfkill unblock all
ExecStart=/usr/bin/aic-bluetooth
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
