# userpatches/overlay/autostart/browser.desktop
[Desktop Entry]
Exec=$HTCDISPLAY_INSTALL_DIR/run-chromium.sh