# userpatches/overlay/autostart/unclutter.desktop
[Desktop Entry]
Exec="nohup unclutter -idle 2 &"