# packages/bsp/orangepi3b/orangepi3b-sprd-bluetooth.service
[Unit]
Description=Spreadtrum (sprd) bluetooth support for Orange Pi 3B v1.1
Before=bluetooth.service

[Service]
Type=simple
RemainAfterExit=Yes
ExecStart=/usr/bin/orangepi3b-sprd-bluetooth

[Install]
WantedBy=multi-user.target
