# packages/bsp/sunxi/ap6255-wifi.service
[Unit]
Description=WIFI AP6255

[Service]
Type=forking
ExecStart=/sbin/modprobe bcmdhd firmware_path=/lib/firmware/fw_bcm43455c0_ag.bin nvram_path=/lib/firmware/nvram_ap6255.txt
TimeoutSec=0
RemainAfterExit=yes
SysVStartPriority=99

[Install]
WantedBy=multi-user.target
