# packages/bsp/sunxi/aw859a-wifi.service
[Unit]
Description=WIFI AW859A

[Service]
Type=forking
ExecStart=/sbin/modprobe sprdwl_ng
TimeoutSec=0
RemainAfterExit=yes
SysVStartPriority=99

[Install]
WantedBy=multi-user.target
