# packages/bsp/sunxi/aw859a-bluetooth.service
[Unit]
Description=Bluetooth AW859A
After=bluetooth.service

[Service]
Type=simple
ExecStartPre=/usr/sbin/rfkill unblock all
ExecStart=/usr/bin/hciattach_opi -s 1500000 /dev/ttyBT0 sprd
TimeoutSec=0
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
