# packages/bsp/pinebook/pinebook-bluetooth.service
[Unit]
Description=Realtek H5 bluetooth support
Before=bluetooth.service

[Service]
ExecStart=/usr/bin/rtk_hciattach -n -s 115200 /dev/ttyS1 rtk_h5

[Install]
WantedBy=multi-user.target
