# config/desktop/_all_distributions/environments/_all_environments/debian/postinst
#
# copy skel to all local users, fix permission and add sudoers for profile sync daemon
#

while IFS=':' read -r login pass uid gid uname homedir comment; do

	if [ "$gid" -ge 1000 ] && [ "$gid" -lt "60000" ]; then

		# copy skel
		sudo cp -R /etc/skel/. $homedir

        	# fix permission
		find /etc/skel/ -mindepth 1 -maxdepth 3 -type d | sed "s|\/etc\/skel|$homedir|" | xargs chown $uid.$gid
        
		# copy PSD
		if ! sudo grep "psd-overlay-helper" /etc/sudoers; then        
			echo "${login} ALL=(ALL) NOPASSWD: /usr/bin/psd-overlay-helper" | sudo tee -a /etc/sudoers
		fi
	fi

done < /etc/passwd

#
# enable profile sync daemon
#

systemctl --user enable psd.service || true >/dev/null 2>&1
systemctl --user start psd.service || true >/dev/null 2>&1

#
# gnome logo hack
#

cp /usr/share/pixmaps/armbian/armbian.png /usr/share/pixmaps/ubuntu-logo-icon.png
cp /usr/share/pixmaps/armbian/armbian.png /usr/share/pixmaps/ubuntu-logo-dark.png
