# config/desktop/_all_distributions/environments/_all_environments/debian/postinst
#
# copy skel to all local users
#
for user_home in $(awk -F'[:]' '{if ($3 >= 1000 && $3 != 65534) print $6}' /etc/passwd)
do
        sudo cp -R /etc/skel/. $user_home
done

#
# fix permission for local users and add sudoers for profile sync daemon
#
while IFS=':' read -r login pass uid gid uname homedir comment; do
        if echo "$homedir" | awk '!/home/{exit 1}'; then
                if [ -d "$homedir" ]; then
                        chown -R $uid:$gid "$homedir";
                        echo "${login} ALL=(ALL) NOPASSWD: /usr/bin/psd-overlay-helper" >> /etc/sudoers
                fi
        fi
done < /etc/passwd

#
# enable profile sync daemon
#

systemctl --user enable psd.service >/dev/null 2>&1
systemctl --user start psd.service >/dev/null 2>&1
