# config/optional/boards/anbernic-rg-ds/_packages/bsp-cli/etc/systemd/system/rgds-pad2key.service
[Unit]
Description=RG DS gamepad-to-keyboard remap
[Service]
ExecStart=/usr/bin/rgds-pad2key
Restart=always
RestartSec=2
[Install]
WantedBy=multi-user.target
