# config/packages/03-common/dedicated/99-armbian-bsp/overlay/lib/systemd/system/armbian-hardware-optimize.service
# Armbian hardware optimization service
# Apply optimisations
# This service may block the boot process for up to 2 minutes

[Unit]
Description=Armbian hardware optimization
Before=basic.target
After=sysinit.target local-fs.target armbian-hardware-monitor.target
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/lib/armbian/armbian-hardware-optimization start
ExecStop=/usr/lib/armbian/armbian-hardware-optimization stop
RemainAfterExit=yes
TimeoutStartSec=2min

[Install]
WantedBy=basic.target
