# packages/bsp/common/etc/systemd/system/openvfd.service
[Unit]
Description=OpenVFD Service
ConditionPathExists=/proc/device-tree/openvfd/

[Service]
ExecStart=/bin/sh -c '[ `cat /proc/device-tree/openvfd/compatible` = "open,vfd" ] && /sbin/modprobe openvfd; /usr/sbin/OpenVFDService'
ExecStop=/bin/kill -TERM $MAINPID
ExecStopPost=-/usr/sbin/rmmod openvfd
RemainAfterExit=yes

[Install]
WantedBy=basic.target

#  created by @paradigman (Armbian forum)
