# config/packages/99-board/dedicated/helios4/overlay/lib/systemd/system/helios4-wol.service
[Unit]
Description=Enable Wake-on-LAN for Helios4 eth0
Requires=network.target
After=network.target

[Service]
ExecStart=/sbin/ethtool -s eth0 wol g
Type=oneshot

[Install]
WantedBy=multi-user.target
