# packages/blobs/phytiumpi/rtlbt/systemd-hciattach.service
[Unit]
Description=rtk_hciattach
DefaultDependencies=no
After=NetworkManager.service

[Service]
ExecStart=/usr/bin/rtk_hciattach -n -s 115200 ttyAMA0 rtk_h5 
Restart=on-failure

[Install]
WantedBy=multi-user.target
