# packages/bsp/helios64/helios64-heartbeat-led-legacy.service
[Unit]
Description=Enable heartbeat led on Helios64
After=getty.target

[Service]
Type=oneshot
ExecStart=bash -c 'echo heartbeat | tee /sys/class/leds/helios64\\:\\:status/trigger'

[Install]
WantedBy=multi-user.target
