# packages/bsp/helios64/helios64-ups.timer
[Unit]
Description=Helios64 UPS Shutdown timer on power loss

[Timer]
OnActiveSec=600
AccuracySec=1s

[Install]
WantedBy=timers.target
