# packages/bsp/pocketbeagle2-usb-gadget/pb2-usb-gadget.service
[Unit]
Description=PocketBeagle 2 USB NCM ethernet plus ACM serial gadget
After=sys-kernel-config.mount systemd-modules-load.service
Wants=sys-kernel-config.mount
Before=network-pre.target NetworkManager.service systemd-networkd.service serial-getty@ttyGS0.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/pb2-usb-gadget-ncm-acm start
ExecStop=/usr/local/sbin/pb2-usb-gadget-ncm-acm stop

[Install]
WantedBy=multi-user.target
