# config/packages/99-board/dedicated/a5xmaxplus/overlay/lib/systemd/system/skykirin-ht1628.service
[Unit]
Description=SKY Kirin-Ht1628 LED timer service

[Service]
Type=simple
ExecStart=-/usr/bin/skykirin-ht1628 -d

[Install]
WantedBy=basic.target
