# packages/blobs/station/firefly-fan.service
[Unit]
Description=start firefly-fan service

[Service]
Type=simple
ExecStart=/usr/bin/firefly-fan-init

[Install]
WantedBy=local-fs.target
