# extensions/network/config-networkd/systemd/system/systemd-networkd-wait-online.service.d/armbian-any-interface.conf
[Service]
ExecStart=
ExecStart=/lib/systemd/systemd-networkd-wait-online --any --timeout=20
