# config/packages/03-common/dedicated/02-armbian-desktop-xfce/overlay/etc/systemd/system/desktop-splash.service
[Unit]
Description=Armbian desktop loading splash
DefaultDependencies=no
Before=rsyslog.service syslog.target systemd-journald.service sysinit.target shutdown.target
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/bin/fbi -vt 1 -noverbose -d /dev/fb0 /boot/boot-desktop.png
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target
