# scripts/armhwinfo
#! /bin/bash
### BEGIN INIT INFO
# Provides:          armhwinfo
# Required-Start:    
# Required-Stop:     glibc
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Get some info about hardware for some ARM boards
### END INIT INFO

CORES=$(grep -c ^processor /proc/cpuinfo)
MEMTOTAL=$(awk 'BEGIN { printf "%.0f\n", '$(grep MemTotal /proc/meminfo | awk '{print $2}')'/1024/1024 }')
ARCH=$(lscpu | grep Architecture  | awk '{print $2}')
RTC=$(dmesg | grep rtc0 | awk '{print $(NF)}')
HB_PCI=$(dmesg | grep '16c3:abcd')
HARDWARE=$(cat /proc/cpuinfo | grep Hardware | awk '{print $3}')
GMAC=$(dmesg | grep "sun6i_gmac")
LEDS=$(dmesg |grep "green:ph02:led1")
TERMINUS=$(lsusb  | grep "1a40:0101")

case "$1" in
    start)
# ARM
if [ "$ARCH" = "armv7l" ]; then
       if [ $HARDWARE = "sun7i" ] || [ $HARDWARE = "Allwinner" ]; then
                ID="Allwinner"
                if [ "$MEMTOTAL" = 2 ]; then
                        ID="Cubietruck"
                        else
                        if [ -n "$GMAC" ]; then
                                ID="Banana PI"
								if [ "$TERMINUS" != "" ]; then
                                        ID="Orange PI"
                                fi
                                # there is more
                                if [ "$LEDS" != "" ]; then
                                        ID="Lime 2"
                                fi
                                else
                                ID="Cubieboard"
                                # there is more
                                if [ "$LEDS" != "" ]; then
                                        ID="Lime"
                                fi
                        fi
                fi
        fi
        if [ "$HARDWARE" = "Freescale" ]; then
                ID="Freescale"
                if [ "$MEMTOTAL" = 2 ]; then
                        ID="Cubox i4"
                        else
                        if [ "$HB_PCI" != "" ]; then
                                ID="HB i2eX"
                                else
                                if [ "$RTC" = "rtc0" ]; then
                                        ID="Cubox i2eX"
                                        else
                                        ID="HB i2"
                                        if [ "$CORES" = 1 ]; then ID="HB i1"; fi
                                fi
                        fi
                fi
        fi
fi

# mainline kernel - read from machine info

MACHINE=$(dmesg |grep "Machine model")

if [[ $MACHINE == *LIME2 ]]; then ID="Lime 2"; fi
if [[ $MACHINE == *LIME ]]; then ID="Lime"; fi
if [[ $MACHINE == *"Banana Pi" ]]; then ID="Banana PI"; fi

echo "" > /var/run/motd.dynamic
toilet -f standard -F metal  "$ID" >> /var/run/motd.dynamic
echo "" >> /var/run/motd.dynamic
;;
    *)
        ## If no parameters are given, print which are avaiable.
        echo "Usage: $0 {start}"
        exit 1
        ;;
esac
