# scripts/firstrun.service
[Unit]
Description=Armbian first run tasks
After=local-fs.target network.target
Before=getty.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/etc/init.d/firstrun start

[Install]
WantedBy=basic.target
