# packages/bsp/odroid/odroid-c2-hdmi.service
[Unit]
Description=Odroid C2 HDMI init
DefaultDependencies=no
Wants=rsyslog.service systemd-journald.service
Before=syslog.target sysinit.target
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/lib/armbian/c2_init.sh
RemainAfterExit=no

[Install]
WantedBy=sysinit.target
