# packages/bsp/rv1106/rv1106-firstboot-makeswap.service
[Unit]
Description=Create a swapfile on first boot (board-specific)
ConditionPathExists=!/swapfile
DefaultDependencies=no
After=armbian-resize-filesystem.service
Before=basic.target

[Service]
Type=oneshot
ExecStart=/usr/bin/rv1106-firstboot-makeswap.sh
RemainAfterExit=yes

[Install]
WantedBy=basic.target
