# packages/bsp/armbian-lowmem/lowmem-mkswap.service
[Unit]
Description=Armbian lowmem: Create a swapfile
DefaultDependencies=no
Wants=armbian-resize-filesystem.service
After=armbian-resize-filesystem.service
Before=basic.target

[Service]
Type=oneshot
ExecStart=/usr/bin/lowmem-mkswap.sh
RemainAfterExit=yes

[Install]
WantedBy=basic.target
