# packages/bsp/photonicat2/usb-hub-watchdog/photonicat-usb-hub-watchdog.service
[Unit]
Description=Photonicat USB hub watchdog
ConditionPathExists=/sys/bus/platform/drivers/onboard-usb-dev/unbind
ConditionPathExists=/sys/bus/platform/drivers/onboard-usb-dev/bind

[Service]
Type=oneshot
ExecStart=/usr/lib/armbian/photonicat2-usb-hub-watchdog-run
TimeoutStartSec=120s
StandardOutput=journal
StandardError=journal
