# packages/bsp/photonicat2/usb-hub-watchdog/photonicat-usb-hub-watchdog.timer
[Unit]
Description=Photonicat USB hub watchdog timer

[Timer]
OnBootSec=20s
OnUnitActiveSec=2min
AccuracySec=1s

[Install]
WantedBy=timers.target
