# packages/bsp/rockpis/etc/udev/rules.d/05-fixMACaddress.rules
#If a network interface is being assigned a different MAC address on each boot,
#or the MAC address is based on a disk image (rather than a hardware serial #),
#enable the corresponding line below to derive that interface's MAC address from
#the RK3308 SOC's unique serial number.

#All the two digit hex code prefixes passed to fixEtherAddr should be unique
#and chosen such that (n-2)%4 == 0

KERNEL=="wlan0", ACTION=="add" RUN+="fixEtherAddr %k 0a"
KERNEL=="p2p0", ACTION=="add" RUN+="fixEtherAddr %k 0e"

#U-Boot >=V2024.10 sets the built-in Ethernet MAC adr from the SOC serial number
#KERNEL=="end0", ACTION=="add" RUN+="fixEtherAddr %k 06"
