# packages/bsp/rockpis/etc/udev/rules.d/05-fixMACaddress.rules
#If a network interface is being assigned a new, different address on each boot,
#enable the corresponding line below to derive its MAC addr from UUID of rootfs
#Beware that all the two digit hex code prefixes below must be unique!

KERNEL=="wlan0", ACTION=="add" RUN+="fixEtherAddr %k 0a"
KERNEL=="p2p0", ACTION=="add" RUN+="fixEtherAddr %k 0e"
