# packages/bsp/jethub/jethubj80/bluetooth/jethub-rtk-hciattach.service

[Unit]
Description=rtk_hciattach
DefaultDependencies=no
After=local-fs.target
Before=armbian-ramlog.target
Conflicts=shutdown.target

[Service]
Type=oneshot
ExecStart=/usr/bin/bash /usr/lib/armbian/jethub-rtk-hciattach-starter
RemainAfterExit=yes
TimeoutStartSec=30sec

[Install]
WantedBy=sysinit.target
