# packages/bsp/mvebu/helios4/helios4-wol.service
[Unit]
Description=Enable Wake-on-LAN for Helios4 eth0
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/sbin/ethtool -s eth0 wol g
Type=oneshot

[Install]
WantedBy=multi-user.target
