# packages/bsp/mvebu/helios4/helios4-wol.service
[Unit]
Description=Enable Wake-on-LAN for Helios4 eth0
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/bin/sh -c '/bin/ls -1 /sys/class/net/ | /bin/grep -E "^(eth|en[do])[0-9]+" | /bin/xargs -n1 -I{} /sbin/ethtool -s "{}" wol g'
Type=oneshot

[Install]
WantedBy=multi-user.target
