# packages/bsp/mvebu/helios4/helios4-mdadm-fault-led.service
[Unit]
Description=Helios4 RAID fault LED indicator
Documentation=man:mdadm(8)
ConditionPathExists=/sbin/mdadm
ConditionPathExists=/sys/class/leds/helios4:red:fault
After=local-fs.target

[Service]
Type=simple
ExecStart=/sbin/mdadm --monitor --program=/usr/sbin/mdadm-fault-led.sh
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
