# packages/bsp/mvebu/helios4/helios4-wol.service
[Unit]
Description=Enable Wake-on-LAN for Helios4 eth0
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/bin/sh -c 'for DEV in `ls /sys/class/net/ | egrep "^(en|eth)"` ; do /sbin/ethtool -s "${DEV:?}" wol g ; done'
Type=oneshot

[Install]
WantedBy=multi-user.target
