# config/cli/bsp/helios64/helios64-ups.service
[Unit]
Description=Helios64 UPS Action

[Service]
Type=oneshot
ExecStart=/usr/bin/helios64-ups.sh
