# packages/bsp/orangepi5/orangepi5-usb2-init.service
[Unit]
Description=Init USB2 for Orange Pi 5

[Service]
ExecStart=/usr/bin/sh -c "echo host > /sys/kernel/debug/usb/fc000000.usb/mode"
Type=oneshot

[Install]
WantedBy=default.target
