# packages/blobs/riscv64/thead/lpi4a-sysfan.service
[Unit]
Description=Lpi4A-SysFan
After=multi-user.target

[Service]
ExecStart=/opt/lpi4a-sysfan.sh
Type=simple

[Install]
WantedBy=multi-user.target
