# packages/bsp/clockworkpi-a06/clockworkpi-a06-fan-daemon.service
[Unit]
Description=ClockworkPi DevTerm A06 fan control daemon

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/share/clockworkpi-a06-fan-daemon/bin/temp_fan_daemon_a06.py

[Install]
WantedBy=multi-user.target
