# packages/bsp/common/lib/systemd/system/armbian-led-state.service
# Armbian led state save and restore
# Stores the current led state at shutdown and restores
# it during bootstrap

[Unit]
Description=Armbian leds state
Before=sysinit.target
After=local-fs.target
Conflicts=shutdown.target
DefaultDependencies=no

[Service]
Type=oneshot
RemainAfterExit=true
ReadWritePaths=/sys/class/leds
ExecStart=/usr/lib/armbian/armbian-led-state-restore.sh
ExecStop=/usr/lib/armbian/armbian-led-state-save.sh

[Install]
WantedBy=basic.target

