# packages/bsp/common/lib/systemd/system/armbian-firstrun.service
# Armbian firstrun service
# Creates the swap file, regenerates the SSH keys,
# performs additional first run tweaks
# This service will run in parallel with other services

[Unit]
Description=Armbian first run tasks
Before=getty.target system-getty.slice

[Service]
Type=simple
RemainAfterExit=yes
ExecStart=/usr/lib/armbian/armbian-firstrun start

[Install]
WantedBy=multi-user.target
