# packages/bsp/common/lib/systemd/system/armbian-disable-autologin.timer
[Unit]
Description=Disable automated desktop login

[Timer]
OnBootSec=4min

[Install]
WantedBy=timers.target
