# packages/bsp/common/lib/systemd/system/armbian-disable-autologin.service
[Unit]
Description=Disable automated desktop login

[Service]
Type=oneshot
ExecStart=-rm -f /etc/lightdm/lightdm.conf.d/22-armbian-autologin.conf
ExecStart=-rm -f /etc/sddm.conf.d/autologin.conf
TimeoutStopSec=10
ExecStop=systemctl disable armbian-disable-autologin.timer
