# packages/bsp/common/lib/systemd/system/armbian-firstrun.service
# Armbian firstrun service
# Regenerates the SSH keys and performs additional first run tweaks
# This service will run in parallel with other services

[Unit]
Description=Armbian first run tasks
Before=getty.target system-getty.slice

[Service]
Type=simple
RemainAfterExit=yes
ExecStart=/usr/lib/armbian/armbian-firstrun start
TimeoutStartSec=2min

[Install]
WantedBy=multi-user.target

