# packages/bsp/common/lib/systemd/system/armbian-zram-config.service
[Unit]
Description=Armbian ZRAM config
DefaultDependencies=no
Before=rsyslog.service syslog.target systemd-journald.service sysinit.target shutdown.target armbian-ram-logging.target
After=local-fs.target
Conflicts=shutdown.target
RequiresMountsFor=/var/log /var/log.hdd
IgnoreOnIsolate=yes

[Service]
Type=oneshot
ExecStart=/usr/lib/armbian/armbian-zram-config start
ExecStop=/usr/lib/armbian/armbian-zram-config stop
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target
