# packages/bsp/common/lib/systemd/system/bootsplash-hide-when-booted.service
[Unit]
Description=Hide bootsplash after starting the system
After=graphical.target
DefaultDependencies=no
ConditionPathExists=/sys/devices/platform/bootsplash.0/enabled

[Service]
ExecStartPre=/bin/sh -c 'exec /bin/sleep 15'
ExecStart=/bin/sh -c 'exec /bin/echo off > /sys/devices/platform/bootsplash.0/enabled'
Restart=no
Type=oneshot

[Install]
WantedBy=graphical.target
