# packages/bsp/common/lib/systemd/system/firstrun-config.service
# Armbian firstrun config service
# This service will run in parallel with other services

[Unit]
Description=Armbian first run optional user configuration
Before=network-pre.target network.target
Wants=NetworkManager.service
After=local-fs.target
ConditionPathExists=/boot/armbian_first_run.txt
# This is not supported on Jessie
#ConditionFirstBoot=true

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash /usr/lib/armbian/firstrun-config.sh

[Install]
WantedBy=multi-user.target
