# packages/bsp/common/lib/systemd/system/firstrun-config.service
# Armbian firstrun config service
# This service will run in parallel with other services

[Unit]
Description=Armbian first run optional user configuration
Wants=network-online.target
After=network.target network-online.target
ConditionPathExists=/boot/armbian_first_run.txt
ConditionPathExists=/root/.not_logged_in_yet

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash /usr/lib/armbian/firstrun-config.sh

[Install]
WantedBy=multi-user.target
