# packages/bsp/jethub/jethub-ethreset.service
[Unit]
Description=test service
Wants=network.target
After=network.target

[Service]
Type=oneshot
ExecStartPre=/usr/bin/sleep 20
ExecStart=/usr/lib/armbian/jethub-ethreset

[Install]
WantedBy=multi-user.target