# packages/bsp/jethub/jethub-initer.service
[Unit]
Description=JetHub hardware initializer
Before=basic.target
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/lib/armbian/jethub-init
TimeoutStartSec=1min

[Install]
WantedBy=basic.target

