# config/desktop/.config/autostart/pasystray.desktop
[Desktop Entry]
Encoding=UTF-8
Version=0.9.4
Type=Application
Name=Pulse audio pasystray
Comment=pasystray
Exec=pasystray
OnlyShowIn=XFCE;
StartupNotify=false
Terminal=false
Hidden=false

