# packages/bsp/generate-bt-mac-addr/bt-fixed-mac.service
[Unit]
Description=Bluetooth Fixed Mac Address
After=bluetooth.service systemd-machine-id-commit.service
Wants=bluetooth.service
[Service]
Type=oneshot
RemainAfterExit=yes
SyslogIdentifier=bt-fixed-mac
ExecStart=/usr/local/bin/bt-fixed-mac.sh
[Install]
WantedBy=bluetooth.target
