# config/packages/50-family/dedicated/rk3399/overlay/lib/systemd/system/nanopi-bluetooth.service
[Unit]
Description=Bluetooth nanopi
After=bluetooth.target

[Service]
Type=forking
ExecStartPre=/bin/echo 1 > /sys/class/rfkill/rfkill0/state
ExecStart=/usr/bin/hciattach /dev/ttyS0 bcm43xx 115200 > /dev/null
TimeoutSec=0
RemainAfterExit=yes
SysVStartPriority=99

[Install]
WantedBy=multi-user.target
